# background_script/com.user.webpconvert.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.webpconvert</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/YOUR_USERNAME/scripts/convert_webp.sh</string>
    </array>

    <key>WatchPaths</key>
    <array>
        <string>/Users/YOUR_USERNAME/Desktop/ToWebP</string>
    </array>

    <key>ThrottleInterval</key>
    <integer>2</integer>

    <key>StandardOutPath</key>
    <string>/tmp/webpconvert.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/webpconvert.err</string>
</dict>
</plist>