# background_script/com.user.simshotsorter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.simshotsorter</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/YOUR_USERNAME/scripts/sim_shot_sorter.sh</string>
    </array>

    <key>WatchPaths</key>
    <array>
        <string>/Users/YOUR_USERNAME/Desktop</string>
    </array>

    <key>ThrottleInterval</key>
    <integer>5</integer>

    <key>StandardOutPath</key>
    <string>/tmp/sim_sorter.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/sim_sorter.err</string>
</dict>
</plist>